# deploy/mbg.service
# systemd unit — passerelle BLE→MQTT sur le RPi.
# Install : cp deploy/mbg.service /etc/systemd/system/ ; adapter les Environment=
#           systemctl daemon-reload && systemctl enable --now mbg
[Unit]
Description=Meshtastic BLE -> MQTT gateway (MeshForge)
After=bluetooth.target network-online.target
Wants=bluetooth.target network-online.target

[Service]
# Type=notify + WatchdogSec : l'app envoie WATCHDOG=1 à chaque cycle sain ; si
# systemd cesse de les recevoir (gel total), il relance le service. Filet ultime
# en plus de la reconnexion interne (sonde de vivacité + backoff).
Type=notify
NotifyAccess=main
WatchdogSec=120
# Utilisateur dédié, membre du groupe `bluetooth` pour accéder à hci0.
User=mbg
SupplementaryGroups=bluetooth
WorkingDirectory=/opt/meshtastic-ble-gateway
# Config par variables d'env (cf. src/mbg/config.py). MAC côté BlueZ.
Environment=MBG_BLE_ADDRESS=E6:E3:53:4B:BE:A5
Environment=MBG_BROKER_HOST=mqtt-mt.meteor-oi.re
Environment=MBG_BROKER_PORT=1883
# Backoff de reconnexion : initial / plafond (optionnel, défauts 5 / 30 s).
# Environment=MBG_RECONNECT_DELAY=5
# Environment=MBG_MAX_RECONNECT_DELAY=30
# API de contrôle (downlink) — décommenter + définir un token pour l'activer.
# SECRET : ne pas commiter. Écouter sur le LAN/VPN ; port à ouvrir au pare-feu.
# Environment=MBG_API_TOKEN=change-me
# Environment=MBG_API_HOST=0.0.0.0
# Environment=MBG_API_PORT=8080
# Monitoring / sonde (métriques node + qualité BLE) — 0 = désactivé.
# Prévoir les droits : install -d -o mbg -g mbg /var/lib/mbg (base + CSV).
# Environment=MBG_MONITOR_INTERVAL=300
# Environment=MBG_DB_PATH=/var/lib/mbg/metrics.db
# Environment=MBG_MONITOR_FORCE_TELEMETRY=false
# Environment=MBG_DUMP_DIR=/var/lib/mbg/csv
# Environment=MBG_DUMP_INTERVAL=3600
# Environment=MBG_RETENTION_DAYS=0
ExecStart=/opt/meshtastic-ble-gateway/.venv/bin/python -m mbg
# BLE 24/7 peu fiable : on relance toujours (le code reconnecte déjà en interne).
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
